# roles/workstation/templates/drive-folder.mount.j2
[Unit]
Description={{ item }} Google Drive folder

[Mount]
What={{ home_dir_path }}/Drive/{{ item }}
Where={{ home_dir_path }}/{{ item }}
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
